# proxy/rootfs/etc/services.d/openresty/run
#!/usr/bin/execlineb -P

/usr/local/openresty/bin/openresty -g "daemon off;"
